# installer/linux/flux-updater.path
[Unit]
Description=Watch for Flux update requests

[Path]
PathExists=/var/lib/flux/update-requested

[Install]
WantedBy=multi-user.target
