# installer/linux/flux.service
[Unit]
Description=Flux UPS Monitor
Documentation=https://github.com/oculus-pllx/Flux-public
After=network.target

[Service]
Type=simple
User=flux
Group=flux
EnvironmentFile=/etc/flux/flux.env
ExecStart=/usr/bin/node /opt/flux/backend/server.js
WorkingDirectory=/opt/flux/backend
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=flux

[Install]
WantedBy=multi-user.target
